# .config/autostart/Combin.desktop
[Desktop Entry]
Type=Application
Name=Combin
Exec=/usr/lib/combin/combin.sh --autostart
X-GNOME-Autostart-enabled=true